# apk/aports/e52c/e52c-openrc-services/e52c-dev-perms.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="e52c-dev-perms"
description="Apply custom device-node permissions for fixed peripherals"

depend() {
  need modules
  before networking
}

start() {
  ebegin "Applying custom device permissions"
  /usr/sbin/e52c-apply-dev-perms
  eend 0
}
