# apk/aports/rpi4/rpi4-openrc-services/rpi4-dev-perms.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="rpi4-dev-perms"
description="Apply custom device-node permissions for fixed peripherals"

depend() {
  need modules
  before networking
}

start() {
  ebegin "Applying custom device permissions"
  /usr/sbin/rpi4-apply-dev-perms
  eend 0
}
