# build/aports/rock5b/rock5b-openrc-services/rock5b-dev-perms.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="rock5b-dev-perms"
description="Apply custom device-node permissions for fixed peripherals"

depend() {
  need modules
  before networking
}

start() {
  ebegin "Applying custom device permissions"
  /usr/sbin/rock5b-apply-dev-perms
  eend 0
}
