# apk/aports/e25/e25-openrc-services/e25-dev-perms.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="e25-dev-perms"
description="Apply custom device-node permissions for fixed peripherals"

depend() {
  need modules
  before networking
}

start() {
  ebegin "Applying custom device permissions"
  /usr/sbin/e25-apply-dev-perms
  eend 0
}
