# apk/aports/e54c/e54c-openrc-services/e54c-dev-perms.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="e54c-dev-perms"
description="Apply custom device-node permissions for fixed peripherals"

depend() {
  need modules
  before networking
}

start() {
  ebegin "Applying custom device permissions"
  /usr/sbin/e54c-apply-dev-perms
  eend 0
}
