# build/aports/e25/e25-openrc-services/e25-net-mac-sort.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="e25-net-mac-sort"
description="Sort Ethernet interfaces by MAC address and assign eth0, eth1, ..."

depend() {
  need modules
  before net
}

start() {
  ebegin "Assigning Ethernet interface names by MAC sort order"
  /usr/sbin/e25-net-mac-sort
  eend 0
}
