# build/aports/rock3b/rock3b-usb-updater-services/rock3b-usb-nvme-update.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="rock3b-usb-nvme-update"
description="Flash bundled NVMe image from USB media, then reboot"

depend() {
  need localmount
  before networking
}

start() {
  ebegin "Running USB to NVMe update flow"
  TARGET_NVME_DEVICE="${target_device:-/dev/nvme0n1}" \
    ROOT_PARTLABEL_REQUIRED="${root_partlabel:-rock3b-updater-rootfs}" \
    TARGET_WAIT_SECONDS="${target_wait_seconds:-120}" \
    /usr/sbin/rock3b-run-usb-update
  eend $?
}
