# apk/aports/e54c/e54c-partition-mount-service/e54c-partition-mount.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

name="e54c-partition-mount"
description="Check a hardware device for partitions with labels, and mount them"
command_background="no"

depend() {
  need devfs
  need localmount
  after modules
}

start() {
  ebegin "Mounting additional storage volumes"
  /usr/sbin/e54c-partition-mount
  eend $?
}
