# apk/aports/alpian/partition-mount-service/partition-mount.initd
#!/sbin/openrc-run
# SPDX-License-Identifier: MIT

description="Check a hardware device for labeled partitions and mount them"
command="/usr/sbin/partition-mount"
command_background="no"

depend() {
  need devfs
  need localmount
  after modules
}
